# etc/rc.d/rc.local
#!/bin/sh
#
# /etc/rc.d/rc.local:  Local system initialization script.
#
# Put any local startup commands in here.  Also, if you have
# anything that needs to be run at shutdown time you can
# make an /etc/rc.d/rc.local_shutdown script and put those
# commands in there.

echo 1 > /proc/sys/vm/swappiness

modprobe tpm_tis force=1 interrupts=0

for DEVICE in TPAD TSCR EHC1 EHC2; do 
  echo $DEVICE > /proc/acpi/wakeup
done
